# deploy/countinghouse.service
[Unit]
Description=Read-side energy cost/accounting service
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/sweeney/countinghouse

[Service]
Type=simple
User=countinghouse
Group=countinghouse

ExecStart=/opt/countinghouse/bin/countinghouse
WorkingDirectory=/var/lib/countinghouse

Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/countinghouse
ReadOnlyPaths=/etc/countinghouse
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=countinghouse

[Install]
WantedBy=multi-user.target
